# profile.d/z99_lmod.sh
PROG_PREFIX=${HOME}/.local/progs
. ${PROG_PREFIX}/lmod/lmod/init/profile
module use ${PROG_PREFIX}/modules
module use ${HOME}/.config/modules
