# profile.d/z99_lmod.sh
MODULEPATH=${MODULEPATH:+${MODULEPATH}:}${HOME}/.stow/modules
MODULEPATH=${MODULEPATH:+${MODULEPATH}:}${HOME}/.config/modules
. ${HOME}/.stow/lmod/lmod/init/profile
